# proj.service
[Unit]
Description=Projektor - serwer usługi do udostępniania ekranu

[Service]
ExecStart=laplace
